# source/n/dehydrated/dehydrated.cron
#!/bin/sh
# This file is installed chmod 644 so that it does not run by default.
# If you've set up dehydrated on this system, then make this file executable
# so that renewals will be done automatically through cron:
# chmod 755 dehydrated
/usr/bin/dehydrated -c > /dev/null
