# source/a/etc/_etc/etc/profile.d/home-profile.d.csh.new
#!/bin/csh

if ( -d "$HOME/.profile.d") then
  foreach i ("$HOME"/.profile.d/*.csh)
    source "$i"
  end
endif

unsetenv i
