# source/t/tetex/profile.d/tetex.csh
#!/bin/csh
# Add path and MANPATH for teTeX:
set path = ( $path /usr/share/texmf/bin )
setenv MANPATH ${MANPATH}:/usr/share/texmf/man
