# source/ap/mc/profile.d/mc.csh
alias mc 'source /usr/share/mc/bin/mc-wrapper.csh'
