# source/ap/mc/profile.d/mc.sh
alias mc='. /usr/share/mc/bin/mc-wrapper.sh'
