# source/a/logrotate/logrotate.cron
#!/bin/sh
/usr/sbin/logrotate /etc/logrotate.conf \
 || { /usr/bin/logger -t logrotate "ALERT - exited abnormally." && false ; }
